# dot_config/systemd/user/atuin.service.tmpl
[Unit]
Description=Atuin daemon

[Service]
Type=simple
ExecStart={{ lookPath "atuin" }} daemon

[Install]
WantedBy=default.target
